# modules/reverse/knowledge/attack/techniques.json
{
  "techniques": [
    {
      "technique_id": "T1059",
      "name": "Command and Scripting Interpreter",
      "tactic": ["Execution"],
      "behavior_categories": ["command_execution"],
      "static_evidence": ["command strings", "script interpreter names", "process creation APIs"],
      "dynamic_evidence": ["process_create telemetry", "command-line telemetry", "parent/child process chain"],
      "analysis_methods": ["strings triage", "import review", "capa rule review", "function-level command builder review"],
      "detection_methods": ["process command-line monitoring", "script host execution analytics", "parent-child process correlation"],
      "recommended_functions": ["strings.extract", "strings.enhanced_extract", "behavior.map_static", "behavior.map_dynamic"],
      "false_positive_notes": ["Installers or admin scripts may contain interpreter strings or launch shell commands."],
      "validation_scenarios": ["controlled_command_line_launch"],
      "validation_samples": ["command_execution_fixture"]
    },
    {
      "technique_id": "T1071",
      "name": "Application Layer Protocol",
      "tactic": ["Command and Control"],
      "behavior_categories": ["network_communication"],
      "static_evidence": ["URLs", "domains", "IP addresses", "WinHTTP or socket APIs"],
      "dynamic_evidence": ["dns_query telemetry", "network_connect telemetry", "http_request or tls_connect telemetry"],
      "analysis_methods": ["IOC extraction", "string review", "network API import review", "external intelligence comparison"],
      "detection_methods": ["network connection telemetry", "DNS logs", "proxy logs", "process-to-destination correlation"],
      "recommended_functions": ["ioc.extract", "strings.enhanced_extract", "behavior.map_static", "behavior.map_dynamic"],
      "false_positive_notes": ["Software update checks and telemetry clients may produce network indicators."],
      "validation_scenarios": ["network_callback_fixture"],
      "validation_samples": ["network_connect_fixture"]
    },
    {
      "technique_id": "T1105",
      "name": "Ingress Tool Transfer",
      "tactic": ["Command and Control"],
      "behavior_categories": ["file_write"],
      "static_evidence": ["download paths", "temporary file paths", "file copy/move/delete APIs", "embedded drop paths"],
      "dynamic_evidence": ["file_create telemetry", "file_write telemetry", "file_move or file_delete telemetry"],
      "analysis_methods": ["string path review", "file API review", "dynamic file telemetry review", "static/dynamic comparison"],
      "detection_methods": ["file creation monitoring", "download and staging path monitoring", "process-to-file write correlation"],
      "recommended_functions": ["pe.deep_parse", "strings.enhanced_extract", "behavior.map_static", "behavior.map_dynamic"],
      "false_positive_notes": ["Many benign programs write files during normal operation; require concrete path and process evidence."],
      "validation_scenarios": ["file_write_fixture"],
      "validation_samples": ["file_write_fixture"]
    },
    {
      "technique_id": "T1547.001",
      "name": "Registry Run Keys / Startup Folder",
      "tactic": ["Persistence", "Privilege Escalation"],
      "behavior_categories": ["registry_persistence"],
      "static_evidence": ["CurrentVersion\\Run strings", "RunOnce strings", "registry set APIs or commands"],
      "dynamic_evidence": ["registry_set telemetry under Run or RunOnce keys", "registry_create telemetry for startup locations"],
      "analysis_methods": ["enhanced string review", "registry API review", "static behavior mapping", "focused registry validation"],
      "detection_methods": ["registry value monitoring", "Sysmon registry event review", "startup persistence analytics"],
      "recommended_functions": ["strings.enhanced_extract", "pe.deep_parse", "behavior.map_static", "behavior.map_dynamic"],
      "false_positive_notes": ["Installers and legitimate updaters may use Run keys."],
      "validation_scenarios": ["registry_run_key_fixture"],
      "validation_samples": ["registry_run_key_fixture"]
    },
    {
      "technique_id": "T1055",
      "name": "Process Injection",
      "tactic": ["Defense Evasion", "Privilege Escalation"],
      "behavior_categories": ["process_injection"],
      "static_evidence": ["OpenProcess/VirtualAllocEx/WriteProcessMemory/CreateRemoteThread APIs", "capa injection candidates"],
      "dynamic_evidence": ["process_access telemetry", "memory_write telemetry", "remote_thread_create telemetry"],
      "analysis_methods": ["capa rule review", "import table review", "function-level API sequence review", "dynamic process access telemetry review"],
      "detection_methods": ["cross-process handle monitoring", "remote thread analytics", "memory write telemetry", "EDR process injection analytics"],
      "recommended_functions": ["tool.capa_analyze", "file.lief_parse", "strings.extract", "behavior.map_static", "behavior.map_dynamic"],
      "false_positive_notes": ["Debuggers, accessibility tools, security tools, and updaters may trigger similar process access evidence."],
      "validation_scenarios": ["controlled_remote_thread_simulation"],
      "validation_samples": ["process_injection_fixture"]
    },
    {
      "technique_id": "T1003",
      "name": "OS Credential Dumping",
      "tactic": ["Credential Access"],
      "behavior_categories": ["credential_access"],
      "static_evidence": ["lsass strings", "credential strings", "token privilege APIs", "mimikatz-like indicators"],
      "dynamic_evidence": ["process access targeting lsass", "command-line references to credential material", "credential-related tool execution"],
      "analysis_methods": ["sensitive string review", "API review", "process access telemetry review", "function-level credential access review"],
      "detection_methods": ["LSASS access monitoring", "credential dumping command detection", "sensitive process access analytics"],
      "recommended_functions": ["strings.enhanced_extract", "tool.capa_analyze", "behavior.map_static", "behavior.map_dynamic"],
      "false_positive_notes": ["Inventory tools, endpoint agents, and diagnostic tools may reference credential terms without dumping secrets."],
      "validation_scenarios": ["credential_access_query_fixture"],
      "validation_samples": ["credential_access_fixture"]
    },
    {
      "technique_id": "T1497",
      "name": "Virtualization/Sandbox Evasion",
      "tactic": ["Defense Evasion", "Discovery"],
      "behavior_categories": ["anti_analysis", "anti_vm_sandbox_evasion"],
      "static_evidence": ["debugger API strings", "VMware/VirtualBox strings", "sandbox or monitoring tool strings"],
      "dynamic_evidence": ["process or registry queries for debugger, sandbox, VM, or analysis tool artifacts"],
      "analysis_methods": ["anti-debug API review", "environment string review", "dynamic command-line and registry query review"],
      "detection_methods": ["analysis-tool discovery telemetry", "VM artifact query monitoring", "anti-debug API telemetry"],
      "recommended_functions": ["strings.enhanced_extract", "behavior.map_static", "behavior.map_dynamic", "tool.ida_function_features"],
      "false_positive_notes": ["Software compatibility checks may query VM or debugger indicators for benign reasons."],
      "validation_scenarios": ["anti_analysis_environment_probe_fixture"],
      "validation_samples": ["anti_analysis_fixture"]
    },
    {
      "technique_id": "T1082",
      "name": "System Information Discovery",
      "tactic": ["Discovery"],
      "behavior_categories": ["system_discovery"],
      "static_evidence": ["GetSystemInfo/GetNativeSystemInfo strings", "computer-name and OS query APIs"],
      "dynamic_evidence": ["API or event evidence for computer name, system information, memory, or OS queries"],
      "analysis_methods": ["string/API review", "Procmon/Noriben API summary review", "sandbox signature comparison"],
      "detection_methods": ["API telemetry", "EDR system discovery analytics", "process behavior correlation"],
      "recommended_functions": ["strings.enhanced_extract", "behavior.map_static", "behavior.map_dynamic"],
      "false_positive_notes": ["Legitimate applications commonly query system information for compatibility."],
      "validation_scenarios": ["system_discovery_fixture"],
      "validation_samples": []
    },
    {
      "technique_id": "T1069",
      "name": "Permission Groups Discovery",
      "tactic": ["Discovery"],
      "behavior_categories": ["privilege_discovery"],
      "static_evidence": ["LookupPrivilegeValue strings", "SeDebugPrivilege strings", "token privilege APIs"],
      "dynamic_evidence": ["LookupPrivilegeValue or token privilege query evidence"],
      "analysis_methods": ["API/string review", "Procmon/Noriben API summary review"],
      "detection_methods": ["privilege query telemetry", "token API monitoring"],
      "recommended_functions": ["strings.enhanced_extract", "behavior.map_static", "behavior.map_dynamic"],
      "false_positive_notes": ["Installers and admin tools may query privileges before legitimate operations."],
      "validation_scenarios": ["privilege_discovery_fixture"],
      "validation_samples": []
    },
    {
      "technique_id": "T1083",
      "name": "File and Directory Discovery",
      "tactic": ["Discovery"],
      "behavior_categories": ["file_directory_discovery"],
      "static_evidence": ["FindFirstFile/FindNextFile strings", "directory enumeration strings"],
      "dynamic_evidence": ["directory_enumerated, file query, or FindFirst/FindNext evidence"],
      "analysis_methods": ["API/string review", "Procmon/Noriben timeline review"],
      "detection_methods": ["file enumeration telemetry", "EDR directory discovery analytics"],
      "recommended_functions": ["strings.enhanced_extract", "behavior.map_static", "behavior.map_dynamic"],
      "false_positive_notes": ["Runtime frameworks and benign applications may enumerate framework or resource directories."],
      "validation_scenarios": ["file_directory_discovery_fixture"],
      "validation_samples": []
    },
    {
      "technique_id": "T1112",
      "name": "Modify Registry",
      "tactic": ["Defense Evasion"],
      "behavior_categories": ["trace_setting_modification"],
      "static_evidence": ["Microsoft Tracing/RASAPI32/RASMANCS registry strings", "RegSetValue evidence"],
      "dynamic_evidence": ["registry writes to Microsoft Tracing or related trace configuration keys"],
      "analysis_methods": ["registry string review", "Procmon/Noriben registry write review", "static/dynamic comparison"],
      "detection_methods": ["registry set telemetry", "Sysmon registry events", "Sigma registry modification rules"],
      "recommended_functions": ["strings.enhanced_extract", "behavior.map_static", "behavior.map_dynamic"],
      "false_positive_notes": ["RAS/network components and crash handlers can touch tracing keys during normal or crash recovery paths."],
      "validation_scenarios": ["trace_setting_modification_fixture"],
      "validation_samples": []
    },
    {
      "technique_id": "T1027",
      "name": "Obfuscated Files or Information",
      "tactic": ["Defense Evasion"],
      "behavior_categories": ["packer_or_obfuscation"],
      "static_evidence": ["high entropy sections", "packed section names", "encoded blobs", "obfuscated strings"],
      "dynamic_evidence": ["runtime string decode evidence", "memory unpacking or dump evidence when collected"],
      "analysis_methods": ["entropy review", "packer detection", "FLOSS string recovery", "function-level decode loop review"],
      "detection_methods": ["packed binary heuristics", "high entropy section detection", "runtime unpacking telemetry where available"],
      "recommended_functions": ["file.byte_stats", "packer.detect_enhanced", "tool.floss_extract", "tool.detect_it_easy"],
      "false_positive_notes": ["Compressed installers and legitimate protectors may look packed or high entropy."],
      "validation_scenarios": ["static_packer_obfuscation_fixture"],
      "validation_samples": ["packer_or_obfuscation_fixture"]
    },
    {
      "technique_id": "T1543.003",
      "name": "Windows Service",
      "tactic": ["Persistence", "Privilege Escalation"],
      "behavior_categories": ["service_creation"],
      "static_evidence": ["CreateService/OpenService/StartService APIs", "sc.exe strings", "service registry paths"],
      "dynamic_evidence": ["service_create telemetry", "service_update telemetry", "service_start telemetry"],
      "analysis_methods": ["service API review", "command-line review", "service telemetry review", "registry service path review"],
      "detection_methods": ["service creation monitoring", "service binary path monitoring", "sc.exe and service-control command analytics"],
      "recommended_functions": ["strings.enhanced_extract", "pe.deep_parse", "behavior.map_static", "behavior.map_dynamic"],
      "false_positive_notes": ["Installers and enterprise management tools commonly create or update services."],
      "validation_scenarios": ["service_creation_fixture"],
      "validation_samples": ["service_creation_fixture"]
    },
    {
      "technique_id": "T1053.005",
      "name": "Scheduled Task",
      "tactic": ["Execution", "Persistence", "Privilege Escalation"],
      "behavior_categories": ["scheduled_task"],
      "static_evidence": ["schtasks strings", "ScheduledTask cmdlets", "task names", "Task Scheduler COM/API usage"],
      "dynamic_evidence": ["task_create telemetry", "task_run telemetry", "process creation linked to task scheduler"],
      "analysis_methods": ["string review", "command-line review", "dynamic task telemetry review", "focused task validation"],
      "detection_methods": ["TaskScheduler logs", "process creation telemetry for schtasks or PowerShell task cmdlets", "new scheduled task monitoring"],
      "recommended_functions": ["strings.enhanced_extract", "behavior.map_static", "behavior.map_dynamic"],
      "false_positive_notes": ["Installers and maintenance tools may create scheduled tasks."],
      "validation_scenarios": ["scheduled_task_fixture"],
      "validation_samples": ["scheduled_task_fixture"]
    }
  ]
}
